# chromium-kiosk/rootfs/etc/services.d/hardware/run
#!/bin/bash
exec 2>&1

echo "================================================================="
echo "  [CHECK-LIST] INITIALISATION DES COMPOSANTS SYSTÈME (UDEV)      "
echo "================================================================="

# 1. Amorce de la détection de la dalle tactile et des périphériques d'entrée
echo -n "  [CHECK] Initialisation de UDEV pour le pavé tactile... "
if udevadm trigger --action=add 2>/dev/null; then
    echo -e "\e[32m[ OK ]\e[0m"
else
    echo -e "\e[33m[ ATTENTION ]\e[0m (Vérifier les privilèges)"
fi

# 2. Validation de l'accès direct aux nœuds graphiques du GPU
echo -n "  [CHECK] Détection de l'accélérateur graphique (DRM)... "
if [ -d /dev/dri ] && [ "$(ls -A /dev/dri 2>/dev/null)" ]; then
    echo -e "\e[32m[ OK ]\e[0m"
    ls -l /dev/dri | sed 's/^/         -> /'
else
    echo -e "\e[31m[ ÉCHEC ]\e[0m (Accès direct GPU manquant)"
fi

echo "================================================================="
echo "  [CHECK-LIST TERMINÉE] PRÉPARATION DU MATÉRIEL RÉUSSIE          "
echo "================================================================="

# Création sécurisée des répertoires pour l'affichage
mkdir -p /run/user/0
chmod 0700 /run/user/0

# Ce script a terminé son travail de vérification initiale
exit 0
