# chromium-kiosk/rootfs/etc/services.d/hardware/run
#!/bin/bash
exec 2>&1

echo "================================================================="
echo "  [CHECK-LIST] INITIALISATION DES COMPOSANTS SYSTÈME (SEATD)     "
echo "================================================================="

# 1. Démarrage de seatd
seatd -g video &
PID_SEATD=$!

TIMEOUT=5
while [ ! -S /run/seatd.sock ] && [ $TIMEOUT -gt 0 ]; do
    sleep 0.5
    let TIMEOUT--
done

if [ -S /run/seatd.sock ]; then
    echo "  [ OK ] Étape 1 : Gestionnaire d'accès seatd opérationnel."
else
    echo "  [ ÉCHEC ] Étape 1 : Impossible de lier la socket seatd."
    exit 1
fi

# 2. Amorce de la détection de l'écran tactile
echo -n "  [CHECK] Initialisation de UDEV pour le pavé tactile... "
if udevadm trigger --action=add 2>/dev/null; then
    echo -e "\e[32m[ OK ]\e[0m"
else
    echo -e "\e[33m[ ATTENTION ]\e[0m"
fi

# 3. Validation de l'accès direct aux pilotes graphiques (GPU Intel)
echo -n "  [CHECK] Détection de l'accélérateur graphique (DRM)... "
if [ -d /dev/dri ] && [ "$(ls -A /dev/dri 2>/dev/null)" ]; then
    echo -e "\e[32m[ OK ]\e[0m"
    ls -l /dev/dri | sed 's/^/         -> /'
else
    echo -e "\e[31m[ ÉCHEC ]\e[0m (Accès direct GPU manquant)"
fi

echo "================================================================="
echo "  [CHECK-LIST TERMINÉE] TOUS LES SYSTÈMES MATÉRIELS SONT PRÊTS   "
echo "================================================================="

mkdir -p /run/xdg-kiosk
chmod 0700 /run/xdg-kiosk

trap "kill $PID_SEATD; exit 0" SIGINT SIGTERM
wait $PID_SEATD
