# chromium-kiosk/rootfs/etc/services.d/udevd/run
#!/bin/bash
exec 2>&1
echo "🔌 Démarrage du daemon udevd..."
exec /sbin/udevd --debug
