# chromium-kiosk/rootfs/etc/services.d/seatd/run
#!/bin/bash
exec 2>&1
echo "🪑 Démarrage du démon seatd..."
exec seatd
